# rename/upstart/rename-milter.conf
description	"Header renaming milter"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid milter
setgid milter

exec /usr/bin/python2 /usr/local/lib/milter/rename
