# rename/upstart/rename-milter.conf
description	"Header renaming milter"

start on starting postfix
stop on stopped postfix

respawn
setuid milter
setgid milter

exec /usr/bin/python2 /usr/local/lib/milter/rename
